# launchd/com.researchman.autoresearch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.researchman.autoresearch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>
      export PATH="$HOME/.local/bin:/usr/local/bin:/opt/homebrew/bin:$PATH";
      cd /Users/tm/Projects/Claude/ResearchMan;
      LOG=$HOME/Library/Logs/researchman-auto.log;
      node scripts/run-if-due.mjs &gt;&gt; $LOG 2&gt;&amp;1 || exit 0;
      # Technology側(techresearch)と同時発火した際のgit競合を防ぐ排他ロック（最大30分待ち）
      LOCK=/tmp/researchman-git.lock; waited=0;
      while ! mkdir $LOCK 2&gt;/dev/null; do sleep 30; waited=$((waited+30)); [ $waited -ge 1800 ] &amp;&amp; { echo "ロック取得タイムアウト: $(date)" &gt;&gt; $LOG; exit 0; }; done;
      trap "rmdir $LOCK 2&gt;/dev/null" EXIT;
      echo "===== Run start: $(date) =====" &gt;&gt; $LOG;
      node scripts/auto-research-cc.mjs &gt;&gt; $LOG 2&gt;&amp;1;
      node scripts/self-heal-thumbnails.mjs &gt;&gt; $LOG 2&gt;&amp;1;
      git add data/cases.json public/thumbnails/ &gt;&gt; $LOG 2&gt;&amp;1;
      if git diff --quiet HEAD; then
        echo "変更なし（新規事例なし）: $(date)" &gt;&gt; $LOG;
      else
        git commit -m "Auto research: $(date '+%Y-%m-%d')" &gt;&gt; $LOG 2&gt;&amp;1;
        if git push &gt;&gt; $LOG 2&gt;&amp;1; then
          if node scripts/verify-deploy.mjs &gt;&gt; $LOG 2&gt;&amp;1; then
            echo "反映まで確認OK: $(date)" &gt;&gt; $LOG;
            node scripts/send-mail.mjs &gt;&gt; $LOG 2&gt;&amp;1;
            node scripts/notify-line.mjs &gt;&gt; $LOG 2&gt;&amp;1;
          else
            echo "push成功だが反映未確認（時間切れ）: $(date)" &gt;&gt; $LOG;
          fi;
        else
          echo "push失敗（pre-push監査で中止の可能性）。コミットはローカル残存。要手動対応: $(date)" &gt;&gt; $LOG;
        fi;
      fi;
      echo "Completed: $(date)" &gt;&gt; $LOG
    </string>
  </array>

  <!-- 毎時起動して run-if-due.mjs が「前回から72時間経過」を判定する。
       スリープ/シャットダウンで実行時刻を逃しても、復帰後1時間以内に必ず走る。 -->
  <key>StartInterval</key>
  <integer>3600</integer>

  <!-- ログイン/再起動直後にも判定を走らせる（復旧時の即時キャッチアップ） -->
  <key>RunAtLoad</key>
  <true/>

  <!-- クラッシュ時に再起動しない -->
  <key>KeepAlive</key>
  <false/>

  <!-- ログ（/tmpは再起動で消えるためLibrary/Logsに置く） -->
  <key>StandardOutPath</key>
  <string>/Users/tm/Library/Logs/researchman-auto.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tm/Library/Logs/researchman-auto-error.log</string>
</dict>
</plist>
