# launchd/com.researchman.tuneup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.researchman.tuneup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>
      export PATH="$HOME/.local/bin:/usr/local/bin:/opt/homebrew/bin:$PATH";
      cd /Users/tm/Projects/Claude/ResearchMan;
      LOG=$HOME/Library/Logs/researchman-tuneup.log;
      [ -f $LOG ] &amp;&amp; [ $(stat -f%z $LOG) -gt 5242880 ] &amp;&amp; mv $LOG $LOG.1;
      node scripts/run-if-due.mjs --state .last-tuneup-run.txt --daily-at 8 --minute 30 --monthly-days 1,15 &gt;&gt; $LOG 2&gt;&amp;1 || exit 0;
      # 収集3ジョブ(cc/tech/ideaseeds)と同じ排他ロックで直列化（git競合防止。最大45分待ち）。
      # ロックが90分超残存していたらkill -9等の残骸とみなし奪取する
      LOCK=/tmp/researchman-git.lock; waited=0;
      while ! mkdir $LOCK 2&gt;/dev/null; do
        if [ -d $LOCK ] &amp;&amp; [ $(( $(date +%s) - $(stat -f %m $LOCK) )) -gt 5400 ]; then
          echo "staleロック奪取: $(date)" &gt;&gt; $LOG; rmdir $LOCK 2&gt;/dev/null; continue;
        fi;
        sleep 30; waited=$((waited+30));
        [ $waited -ge 2700 ] &amp;&amp; { echo "ロック取得タイムアウト: $(date)" &gt;&gt; $LOG; exit 0; };
      done;
      trap "rmdir $LOCK 2&gt;/dev/null" EXIT;
      echo "===== Tuneup start: $(date) =====" &gt;&gt; $LOG;
      if node scripts/biweekly-tuneup.mjs &gt;&gt; $LOG 2&gt;&amp;1; then
        git add data/research-tuning.json data/idea-tuning.json data/x-radar-queries.json RESEARCH_PLAN.md &gt;&gt; $LOG 2&gt;&amp;1;
        if git diff --cached --quiet; then
          echo "変更なし（分析結果は現状維持 or スキップ）: $(date)" &gt;&gt; $LOG;
          node scripts/notify-line.mjs --text-file /tmp/researchman-tuneup-report.txt &gt;&gt; $LOG 2&gt;&amp;1;
        else
          git commit -m "chore: biweekly research tuning $(date '+%Y-%m-%d')" &gt;&gt; $LOG 2&gt;&amp;1;
          if git push &gt;&gt; $LOG 2&gt;&amp;1; then
            if node scripts/verify-deploy.mjs --skip-pages &gt;&gt; $LOG 2&gt;&amp;1; then
              echo "反映まで確認OK: $(date)" &gt;&gt; $LOG;
            else
              echo "push成功だが反映未確認（時間切れ）: $(date)" &gt;&gt; $LOG;
            fi;
            node scripts/notify-line.mjs --text-file /tmp/researchman-tuneup-report.txt &gt;&gt; $LOG 2&gt;&amp;1;
          else
            echo "push失敗（pre-push監査で中止の可能性）。コミットはローカル残存。要手動対応: $(date)" &gt;&gt; $LOG;
            node scripts/notify-line.mjs --text-file /tmp/researchman-tuneup-report.txt &gt;&gt; $LOG 2&gt;&amp;1;
          fi;
        fi;
      else
        echo "チューンアップエラー終了: $(date)" &gt;&gt; $LOG;
        node scripts/run-if-due.mjs --state .last-tuneup-run.txt --mark &gt;&gt; $LOG 2&gt;&amp;1;
        node scripts/notify-line.mjs --text-file /tmp/researchman-tuneup-report.txt &gt;&gt; $LOG 2&gt;&amp;1;
      fi;
      echo "Tuneup completed: $(date)" &gt;&gt; $LOG
    </string>
  </array>

  <!-- 毎月1日・15日 08:30に実行。--monthly-days ゲートが対象日以外は即exit 3（ログなし）。
       対象日に08:30時点でPCが落ちていても、以降23:30まで毎正時30分に起動してキャッチアップする。
       他の日次3ジョブと同じ /tmp/researchman-git.lock で排他制御 -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/tm/Library/Logs/researchman-tuneup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tm/Library/Logs/researchman-tuneup-error.log</string>
</dict>
</plist>
