# launchd/com.researchman.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.researchman.watchdog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>
      export PATH="$HOME/.local/bin:/usr/local/bin:/opt/homebrew/bin:$PATH";
      cd /Users/tm/Projects/Claude/ResearchMan;
      LOG=$HOME/Library/Logs/researchman-watchdog.log;
      [ -f $LOG ] &amp;&amp; [ $(stat -f%z $LOG) -gt 5242880 ] &amp;&amp; mv $LOG $LOG.1;
      STATE=.last-watchdog-run.txt;
      AMDUE=0;
      node scripts/run-if-due.mjs --state $STATE --daily-at 12 --minute 30 &gt;&gt; $LOG 2&gt;&amp;1 &amp;&amp; AMDUE=1;
      PMDUE=0;
      if [ "$AMDUE" = "0" ]; then
        node scripts/run-if-due.mjs --state $STATE --daily-at 18 --minute 30 &gt;&gt; $LOG 2&gt;&amp;1 &amp;&amp; PMDUE=1;
      fi;
      if [ "$AMDUE" = "0" ] &amp;&amp; [ "$PMDUE" = "0" ]; then exit 0; fi;
      DEEP="";
      if [ "$PMDUE" = "1" ] &amp;&amp; [ "$(date +%u)" = "7" ]; then DEEP="--deep"; fi;
      echo "===== Watchdog start ($([ "$AMDUE" = "1" ] &amp;&amp; echo AM || echo PM)$([ -n "$DEEP" ] &amp;&amp; echo " +deep")): $(date) =====" &gt;&gt; $LOG;
      node scripts/watchdog.mjs $DEEP &gt;&gt; $LOG 2&gt;&amp;1 || echo "watchdogエラー終了: $(date)" &gt;&gt; $LOG;
      echo "Watchdog completed: $(date)" &gt;&gt; $LOG
    </string>
  </array>

  <!-- 毎日12:30/18:30を毎正時30分でキャッチアップ判定（.last-watchdog-run.txtを共有した
       2段ゲート）。日曜のPM枠のみ--deepを付与。既存3ジョブと違い、このplistは
       /tmp/researchman-git.lockを一切取得しない（watchdog.mjs内のgit操作は
       scripts/lib/watchdog-git.mjsが必要な瞬間だけ短く取る設計。理由: watchdogが
       launchctl kickstartで他ジョブを起こす際、ロックを握ったままだと起こした側が
       ロックを取れず自己デッドロックする実バグを踏んだため、2026-07-08に変更） -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/tm/Library/Logs/researchman-watchdog.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tm/Library/Logs/researchman-watchdog-error.log</string>
</dict>
</plist>
